# systemd/greyline.service
[Unit]
Description=Render the greyline world-time wallpaper
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
# Adjust the path if not installed via pipx (~/.local/bin) — e.g. /usr/bin.
ExecStart=%h/.local/bin/greyline

[Install]
WantedBy=graphical-session.target
