# systemd/greyline.timer
[Unit]
Description=Update the greyline world-time wallpaper on a schedule

[Timer]
OnCalendar=*:*:00
Persistent=true

[Install]
WantedBy=timers.target
