# systemd/thinkpad-wallpaper.timer
[Unit]
Description=Update the ThinkPad World Time wallpaper on a schedule

[Timer]
OnCalendar=*:*:00
Persistent=true

[Install]
WantedBy=timers.target
